# configs/com.opencode.telegram-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.opencode.telegram-bridge</string>

	<key>Program</key>
	<string>/Users/YOUR_USERNAME/opencode-telegram/opencode-telegram</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>TELEGRAM_BOT_TOKEN</key>
		<string>YOUR_BOT_TOKEN_HERE</string>
		<key>TELEGRAM_CHAT_ID</key>
		<string>YOUR_CHAT_ID_HERE</string>
		<key>OPENCODE_BASE_URL</key>
		<string>http://localhost:54321</string>
		<key>OPENCODE_DIRECTORY</key>
		<string>/Users/YOUR_USERNAME/.config/opencode</string>
	</dict>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>

	<key>StandardOutPath</key>
	<string>/Users/YOUR_USERNAME/Library/Logs/opencode-telegram/stdout.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/YOUR_USERNAME/Library/Logs/opencode-telegram/stderr.log</string>

	<key>WorkingDirectory</key>
	<string>/Users/YOUR_USERNAME/opencode-telegram</string>
</dict>
</plist>
